# gui-session.target
[Unit]
Description=GUI Session
After=gui-session-pre.target
#Wants=xdg-desktop-autostart.target
RefuseManualStart=yes
