# gui-session-pre.target
[Unit]
Description=GUI Session Pre
RefuseManualStart=yes

